# rootfs/etc/s6/.s6-svscan/finish
#!/bin/execlineb -S0
echo Finishing... $@
