# rootfs/etc/services.d/docker-gen/run
#!/usr/bin/execlineb -P
fdmove -c 2 1
with-contenv /usr/local/bin/docker-gen -watch -only-exposed -notify "/bin/s6-svc -t /var/run/s6/services/nginx" /etc/docker-gen/nginx.tmpl /etc/nginx/conf.d/default.conf
